# deploy/playbooks/templates/[[app_name]].service.j2
[Unit]
Description={{ app_name }} API (Bun/Hono)
After=network.target

[Service]
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ server_current }}
EnvironmentFile=/etc/{{ app_name }}.env
ExecStart=/root/.bun/bin/bun ./index.ts
Environment=NODE_ENV=production
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
